# users/Profpatsch/whatcd-resolver/services/reverse-proxy/run
#!/usr/bin/env execlineb
caddy reverse-proxy --from :9092 --to :9093
