# third_party/nix/misc/systemd/nix-daemon.service.in
[Unit]
Description=Nix Daemon
RequiresMountsFor=@storedir@
RequiresMountsFor=@localstatedir@
ConditionPathIsReadWrite=@localstatedir@/nix/daemon-socket

[Service]
ExecStart=@@bindir@/nix-daemon nix-daemon
KillMode=process
Environment=GRPC_VERBOSITY=DEBUG
Environment=GRPC_TRACE=server_channel,executor,channel,api,fd_trace,polling

[Install]
WantedBy=multi-user.target
